# x11/wallpaper.systemd/wallpaper.service
[Unit]
Description=User service for wallpaper 

[Service]
Type=oneshot
ExecStart=/bin/bash '$(which wallpaper)'

[Install]
WantedBy=default.target
