# x11/wallpaper.systemd/wallpaper.timer
[Unit]
Description=User Timer for wallpaper 

[Timer]
OnUnitActiveSec=15m
Unit=wallpaper.service

[Install]
WantedBy=default.target
